# kingfisher_bringup/upstart/kingfisher.conf
description	"bringup kingfisher"

start on net-device-up IFACE=wlan0
stop on net-device-down IFACE=wlan0

console output
respawn

exec kingfisher-start
pre-stop exec kingfisher-stop
